# systemd/ai-vault.service
[Unit]
Description=AI Vault - AI-powered knowledge vault
After=network.target

[Service]
Type=simple
User=ai-vault
WorkingDirectory=/opt/ai-vault
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/ai-vault/.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
